# debsetup/appserver/init/networkoptix-entcontroller.conf
# Network Optix Media Server

description "Network Optix Enterprise Controller"
author "Network Optix <support@networkoptix.com>"

# When to start the service
start on (net-device-up
          and local-filesystems
      and runlevel [2345])

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Start the process
exec /opt/networkoptix/entcontroller/bin/entcontroller

