# debsetup/appserver-deb/maven/filter-resources/init/networkoptix-entcontroller.conf
# ${company.name} Media Server

description "${company.name} Enterprise Controller"
author "${company.name} <support@${company.url}>"

kill timeout 120

# When to start the service
start on (net-device-up
          and local-filesystems
      and runlevel [2345])

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Start the process
exec /opt/${deb.customization.company.name}/entcontroller/bin/entcontroller

