# nx_cloud_deploy/monitoring_simple/systemd/vanyas_tests.service
[Unit]
Description=Runs cloud portal tests. Reports to Amazon Cloud Watch.

[Service]
Type=oneshot
ExecStart=/bin/bash -c /etc/systemd/system/vanyas_tests.sh
